# configs/systemd/backup-configs.timer
[Unit]
Description=Daily backup of critical configs (srv-core)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
